# etc/systemd/system/t2bouncer.timer
[Unit]
Description=t2bouncer timer
BindsTo=t2server.service

[Timer]
Unit=t2bouncer.service
OnCalendar=*:00:00

[Install]
WantedBy=timers.target